# roles/insight/templates/insight.service.j2
[Unit]
Description=Insight Workbench Service
After=network.target

[Service]
Type=simple
User={{ workbench_user }}
Group={{ workbench_group }}
WorkingDirectory={{ workbench_dir }}
Environment=PORT={{ workbench_port }}
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
